# cap.crontab
*/1 8-20 * * 1-5   $HOME/Desktop/caps/capit
